# systemd/aimarket-agent4-bot.service
[Unit]
Description=AImarket Agent 4 - Telegram lead qualifier bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/aimarket-agents
EnvironmentFile=/opt/aimarket-agents/.env
ExecStart=/opt/aimarket-agents/venv/bin/python -m agent4_lead_qualifier.bot
Restart=on-failure
RestartSec=5
User=aimarket

[Install]
WantedBy=multi-user.target
